# tmrc.service
[Unit]
Description=Performance Moniter Server

[Service]
ExecStart=/opt/tmrcx/start

[Install]
WantedBy=multi-user.target